# 1. code tunnel/yuze-code-tunnel.service
[Service]
User=yuze
Group=security

WorkingDirectory=/home/yuze/
ExecStart=/home/yuze/code tunnel
Restart=always
RestartSec=10
